# telebot.service
[Unit]
Description=Telegram Bot - Πειραιεύς
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telebot
ExecStart=/home/ubuntu/telebot/venv/bin/python bot.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/telebot/.env

[Install]
WantedBy=multi-user.target
